# docs/infra/supervisor/corex-worker-p24.conf
; ============================================================================
; CoreX — P24 parallel importer queue workers
; ============================================================================
;
; DEPLOY ARTIFACT. Append these two program blocks to the existing
; /etc/supervisor/conf.d/corex-worker.conf on each host that runs the importer
; (live: /corex, demo: /mnt/HC_Volume_103099143/corex-demo), then:
;
;     sudo supervisorctl reread
;     sudo supervisorctl update
;     sudo supervisorctl start corex-worker-p24import:* corex-worker-p24images:*
;     sudo supervisorctl status | grep p24
;
; WHY THIS FILE EXISTS / THE TRAP:
; The existing workers drain ONLY the `default` queue. The parallel importer
; puts confirm jobs on `p24import` and image jobs on `p24images`. Until a worker
; drains each of those, every import job sits in the `jobs` table forever and
; NOTHING happens when a user clicks Import All — no error, just silence. These
; two blocks are not optional; they are the other half of the feature.
;
; THE TWO LANES ARE DELIBERATELY ASYMMETRIC:
;   - p24import  is WIDE (8 procs). A confirm is a DB write, no CDN call, so it
;                can fan out hard. This is where the speed comes from.
;   - p24images  is NARROW (4 procs). Each image job fetches 10-wide, so 4 procs
;                ≈ 40 concurrent requests to images.prop24.com — kept well under
;                the ~80 that tripped P24's per-IP rate limiter and silently lost
;                ~11% of the first import's photos (audit 2026-07-17). DO NOT
;                raise numprocs here to "go faster" — that is exactly the storm
;                the self-healing job + this cap exist to prevent. Speed is the
;                wide property lane; the image lane trades latency for no loss.
;
; The image job carries tries=5 with backoff itself, so a rate-limited gallery
; heals on its own retries rather than sticking short.

[program:corex-worker-p24import]
process_name=%(program_name)s_%(process_num)02d
command=php /corex/artisan queue:work database --queue=p24import --sleep=1 --tries=3 --max-time=3600
directory=/corex
autostart=true
autorestart=true
stopwaitsecs=3600
killasgroup=true
numprocs=8
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/corex/worker-p24import.log
stdout_logfile_maxbytes=20MB

[program:corex-worker-p24images]
process_name=%(program_name)s_%(process_num)02d
command=php /corex/artisan queue:work database --queue=p24images --sleep=1 --tries=5 --max-time=3600 --timeout=310
directory=/corex
autostart=true
autorestart=true
stopwaitsecs=3600
killasgroup=true
numprocs=4
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/corex/worker-p24images.log
stdout_logfile_maxbytes=20MB

; Demo host: same blocks with directory + command paths pointed at
; /mnt/HC_Volume_103099143/corex-demo instead of /corex.
; NOTE: --timeout on p24images (310) must stay just above the job's own
; $timeout (300); a worker --timeout below the job timeout SIGKILLs mid-download.
